# software/usb_ir/packaging/debian/postrm
#!/bin/sh
set -e
update-rc.d iguanaIR remove ||true
rmdir /var/run/iguanaIR 2>/dev/null||true
userdel iguanair||true
